# cron.d/stromding_plotPngs
# MIN   HOUR   MDAY  MON  DOW   COMMAND 
*/5 * * * * root /var/stromding/stromding_plotPngs.pl > /var/log/stromding.log

