# cron.d/stromding_poll2mysql
# MIN   HOUR   MDAY  MON  DOW   COMMAND 
*/1 * * * * root perl /var/stromding/stromding_poll.pl > /var/log/stromding.log
